# recipes/galileo/etc/systemd/system/boot.automount
[Unit]
Before=local-fs.target

[Automount]
Where=/boot
